# systemd/bar.service
[Unit]
Description=Bottom bar

[Service]
ExecStart=/home/hugo/code/c/bottom-bar/bar
Restart=always

[Install]
WantedBy=default.target
